# mongodb.conf
[program:mongodb]
command=/usr/bin/mongo
directory=/var/www/bbs
autostart=true
autorestart=true
